# host-services/coder-agents-sidecars/s6/cliproxy-sidecar/run
#!/command/with-contenv bash
# CLIProxyAPI serving Codex (ChatGPT Plus/Pro OAuth), Gemini (personal Google
# OAuth), and optionally Claude (Anthropic Pro/Max OAuth). Single Go binary,
# single process; the dispatcher routes /v1/messages to it when the model is
# prefixed `subscription/`.
set -euo pipefail

SECRETS_FILE=/run/coder-agents-sidecars/secrets.env
if [ -r "${SECRETS_FILE}" ]; then
  # shellcheck disable=SC1090
  source "${SECRETS_FILE}"
fi

CONFIG_TEMPLATE=/etc/coder-agents-sidecars/cli-proxy.yaml
CONFIG_RUNTIME=/run/coder-agents-sidecars/cli-proxy.yaml
mkdir -p "$(dirname "${CONFIG_RUNTIME}")" /data/auth/cliproxy

# Substitute the shared API key into the runtime config. Done at startup so
# the key can be rotated without rebuilding the image.
KEY="${SIDECAR_SHARED_API_KEY:-}"
if [ -z "${KEY}" ]; then
  echo "[cliproxy-sidecar] SIDECAR_SHARED_API_KEY not loaded; refusing to start" >&2
  exit 1
fi
sed "s|__SIDECAR_SHARED_API_KEY__|${KEY}|g" "${CONFIG_TEMPLATE}" > "${CONFIG_RUNTIME}"

if [ -z "$(ls -A /data/auth/cliproxy 2>/dev/null)" ]; then
  echo "[cliproxy-sidecar] no credential files in /data/auth/cliproxy" >&2
  echo "[cliproxy-sidecar] bootstrap with one or more of:" >&2
  echo "  docker exec -it <container> cli-proxy-api --auth-dir /data/auth/cliproxy --login codex" >&2
  echo "  docker exec -it <container> cli-proxy-api --auth-dir /data/auth/cliproxy --login gemini" >&2
  echo "  docker exec -it <container> cli-proxy-api --auth-dir /data/auth/cliproxy --login claude" >&2
fi

exec cli-proxy-api --config "${CONFIG_RUNTIME}"
