# host-services/coder-agents-sidecars/s6/bootstrap-secrets/up
/opt/headroom-venv/bin/python /usr/local/lib/coder-agents-sidecars/bootstrap_secrets.py
