# host-services/coder-agents-sidecars/s6/gemini-sidecar/run
#!/command/with-contenv bash
# CLIProxyAPI in Gemini-only mode. Exposes native generativelanguage v1beta
# and Cloud Code Assist v1internal endpoints, backed by Gemini CLI OAuth.
set -euo pipefail

CONFIG_TEMPLATE=/etc/coder-agents-sidecars/cli-proxy.yaml
CONFIG_RUNTIME=/run/coder-agents-sidecars/cli-proxy.yaml
mkdir -p "$(dirname "${CONFIG_RUNTIME}")"

# Substitute the shared API key into the runtime config. We do this at startup
# rather than baking into the image so the key can be rotated without rebuild.
KEY="${SIDECAR_SHARED_API_KEY:-changeme}"
sed "s|__SIDECAR_SHARED_API_KEY__|${KEY}|g" "${CONFIG_TEMPLATE}" > "${CONFIG_RUNTIME}"

if [ -z "$(ls -A /data/auth/gemini 2>/dev/null)" ]; then
  echo "[gemini-sidecar] no credential files in /data/auth/gemini" >&2
  echo "[gemini-sidecar] run: docker exec -it <container> cli-proxy-api --auth-dir /data/auth/gemini --login gemini" >&2
fi

exec cli-proxy-api --config "${CONFIG_RUNTIME}"
