# host-services/coder-agents-sidecars/s6/kirocc/run
#!/command/with-contenv bash
# kirocc: Anthropic Messages API → Kiro Builder ID OAuth proxy.
# Reads OAuth state from a Kiro CLI SQLite database; bootstrap by running
# `kiro-auth-login` inside this container.
set -euo pipefail

SECRETS_FILE=/run/coder-agents-sidecars/secrets.env
if [ -r "${SECRETS_FILE}" ]; then
  # shellcheck disable=SC1090
  source "${SECRETS_FILE}"
fi

KIROCC_DATA_DIR="${KIROCC_DATA_DIR:-/data/auth/kirocc}"
KIROCC_DB="${KIROCC_DATA_DIR}/data.sqlite3"
mkdir -p "${KIROCC_DATA_DIR}"

if [ ! -s "${KIROCC_DB}" ]; then
  echo "[kirocc] no SQLite DB at ${KIROCC_DB}" >&2
  echo "[kirocc] bootstrap with:" >&2
  echo "[kirocc]   docker exec -it coder-agents-sidecars kiro-auth-login" >&2
fi

# kirocc collides with Meridian on default port 3456; pin to 9090. Shared API
# key is enforced by kirocc itself via -api-key.
exec kirocc \
  -host 127.0.0.1 \
  -port 9090 \
  -db "${KIROCC_DB}" \
  -api-key "${SIDECAR_SHARED_API_KEY:-}"
