# host-services/coder-agents-sidecars/s6/codex-sidecar/run
#!/command/with-contenv bash
# codex-bridge. Forwards OpenAI Responses API to chatgpt.com/backend-api/codex
# using ChatGPT Plus/Pro Codex CLI auth. Reads CODEX_HOME/auth.json.
set -euo pipefail

if [ ! -f "${CODEX_HOME:-/data/auth/codex}/auth.json" ]; then
  echo "[codex-sidecar] no auth.json at ${CODEX_HOME:-/data/auth/codex}" >&2
  echo "[codex-sidecar] run: docker exec -it <container> codex login" >&2
fi

cd /opt/codex-bridge
exec node bridge.mjs
